# ansible/roles/postgres/templates/postgres-backup.service.j2
[Unit]
Description=PARA PostgreSQL Daily Backup
After=docker.service

[Service]
Type=oneshot
ExecStart={{ para_deploy_dir }}/postgres/scripts/backup.sh
StandardOutput=journal
StandardError=journal
