# ansible/roles/postgres/templates/postgres-backup.timer.j2
[Unit]
Description=PARA PostgreSQL Daily Backup Timer

[Timer]
OnCalendar=daily
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
